# etc/supervisord.d/redis-commander.ini
[program:redis-commander]
command=redis-commander --redis-host redis.service.consul --redis-port 6379
autostart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
